# systemd/aura-ui.service
[Unit]
Description=AURA UI
After=multi-user.target aura.service
Requires=aura.service

[Service]
User=aura
WorkingDirectory=/home/aura/aura
# Wait up to 15 s for the backend socket before launching the UI
ExecStartPre=/bin/bash -c 'for i in $(seq 1 15); do [ -S /tmp/aura.sock ] && exit 0; sleep 1; done; echo "aura.sock not found"; exit 1'
ExecStart=/home/aura/aura/launch_ui.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
