# systemd/aura.service
[Unit]
Description=Aura AI Companion
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=aura
WorkingDirectory=/home/aura/aura
Environment="PATH=/home/aura/aura/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/aura/aura/venv/bin/python3 /home/aura/aura/aura.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aura

[Install]
WantedBy=multi-user.target
