# systemd/aura-csam-log.service
[Unit]
Description=AURA CSAM Privileged Logger
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /home/aura/aura/csam_logger.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
